# deploy/ctools.service
[Unit]
Description=CTools Web Application
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ctools
Environment="PATH=/opt/ctools/venv/bin"
ExecStart=/opt/ctools/venv/bin/python app.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
